# rootfs/etc/s6-overlay/s6-rc.d/unoserver/run
#!/command/with-contenv sh

if [[ -z "${UNOSERVER_CMD}" ]]; then
    UNOSERVER_CMD=unoserver
else
    echo "INFO:unoserver:Start using custom server command"
fi

s6-notifyoncheck $UNOSERVER_CMD
